# server/snoopy.service
[Unit]
Description=Snoopy - Continuous screen recording with web streaming
After=graphical-session.target
Wants=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/bin/snoopy -segment 30m -port 8900
Restart=always
RestartSec=2

[Install]
WantedBy=default.target
