# server/snoopy.service
[Unit]
Description=Snoopy - Continuous screen recording with web streaming
After=graphical-session.target
Wants=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/bin/snoopy -segment 30m -port 8900
# Restart on failure (non-zero exit code) such as when DBus connection is lost
# This allows the service to recover from session pause/lock events
Restart=on-failure
RestartSec=2

[Install]
WantedBy=default.target
